# foreman/master.conf.erb
##
## Warning - this file is generated by mytardis / foreman.  Do not edit it.
##

# (The 'started ssh' event is intended to cater for Ubuntu.) 
start on started sshd or started ssh

stop on starting shutdown

pre-start script

bash << "EOF"
  mkdir -p <%= log_root %>
  chown -R <%= user %> <%= log_root %>
EOF

end script
