# foreman/process_master.conf.erb
##
## Warning - this file is generated by mytardis / foreman.  Do not edit it.
##

start on starting <%= app %>
stop on stopping <%= app %>
